# poshc2.service
[Unit]
Description=PoshC2 Server

[Service]
Type=simple
User=root
ExecStart=/usr/bin/posh-server