# poshc2-docker.service
[Unit]
Description=PoshC2 Dockerised Server

[Service]
Type=simple
User=root
ExecStart=/usr/bin/posh-docker-server

[Install]
WantedBy=default.target