# resources/scripts/poshc2.api.service
[Unit]
Description=PoshC2 API Server

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/posh-api-server

[Install]
WantedBy=default.target