# resources/scripts/poshc2.service
[Unit]
Description=PoshC2 Server

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/posh-server

[Install]
WantedBy=default.target